# worker-agent/worker-agent.service
[Unit]
Description=Movistar Bot Worker Agent
Documentation=https://github.com/tu-usuario/movistar-bot
After=network.target
Wants=network.target

[Service]
Type=simple
User=bot-worker
Group=bot-worker
WorkingDirectory=/opt/worker-agent
ExecStart=/usr/bin/node /opt/worker-agent/dist/index.js
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=movistar-worker

# Security
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/worker-agent/evidence /opt/worker-agent/logs

# Environment
Environment="NODE_ENV=production"
EnvironmentFile=/opt/worker-agent/.env

[Install]
WantedBy=multi-user.target
